# larch/profiles/myarch/rootoverlay/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# Set up automatically logged in user (larch live system only)
if [ -f /.livesys/autologin ]; then
    cp /.livesys/autologin /tmp/newuser
fi

# Restore saved sound volume, etc.
alsactl restore
